# com.groundswell.engine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.groundswell.engine</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/bradwood/Projects/groundswell/run-daemon.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/bradwood/Projects/groundswell</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/bradwood/Projects/groundswell/data/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bradwood/Projects/groundswell/data/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Library/Developer/CommandLineTools/usr/bin</string>
        <key>HOME</key>
        <string>/Users/bradwood</string>
    </dict>
</dict>
</plist>
